# debian/overrides/alsa
### BEGIN INIT INFO
# Provides:          alsa
# Required-Start:    mountall
# Required-Stop:     mountall
# Default-Start:     S
# Default-Stop:      0 1 6
### END INIT INFO
# Reported to Debian BTS as #337638, updated deps after reporting
