# debian/overrides/alsa-utils
### BEGIN INIT INFO
# Provides:          alsa-utils
# Required-Start:    mountall
# Required-Stop:     mountall
# Default-Start:     S
# Default-Stop:      0 6
### END INIT INFO
