# debian/overrides/apmd
### BEGIN INIT INFO
# Provides:          apmd
# Required-Start:    mountvirtfs $syslog
# Required-Stop:     mountvirtfs $syslog
# Default-Start:     S 1 2 3 4 5
# Default-Stop:      0 6
### END INIT INFO
