# debian/overrides/atd
### BEGIN INIT INFO
# Provides:          atd
# Required-Start:    $syslog $time $local_fs $remote_fs
# Required-Stop:     $syslog $time $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO
# Reported to BTS as #335307
