# debian/overrides/bastille-firewall
### BEGIN INIT INFO
# Provides:          bastille-firewall
# Required-Start:    $local_fs ifupdown
# Required-Stop:     $local_fs ifupdown
# Default-Start:     S
# Default-Stop:      0 6
# Short-Description: Load/unload ipchains rulesets
### END INIT INFO
# Reported to BTS as #?
