# debian/overrides/bluez-utils
### BEGIN INIT INFO
# Provides:          bluez-utils
# Required-Start:    $local_fs
# Required-Stop:     $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO
