# debian/overrides/bootlogd
### BEGIN INIT INFO
# Provides:          bootlogd
# Required-Start:    mountdevsubfs
# Required-Stop:
# Default-Start:     S
# Default-Stop:
### END INIT INFO
# Already present in package
