# debian/overrides/bootclean.sh
### BEGIN INIT INFO
# Provides:          bootclean
# Required-Start:    $network
# Required-Stop:     $network
# Default-Start:     1 2 3 4 5
# Default-Stop:      0 6
### END INIT INFO
# Already present in package
