# debian/overrides/checkroot.sh
### BEGIN INIT INFO
# Provides:          checkroot
# Required-Start:    mountvirtfs
# Required-Stop:     
# Should-Start:      keymap hwclockfirst
# Should-stop:
# Default-Start:     S
# Default-Stop:
### END INIT INFO
