# debian/overrides/checkfs.sh
### BEGIN INIT INFO
# Provides:          checkfs
# Required-Start:    modutils checkroot
# Required-Stop: 
# Should-Start:      lvm cryptdisks
# Should-Stop:       lvm cryptdisks
# Default-Start:     S
# Default-Stop:
### END INIT INFO
