# debian/overrides/check-network-interfaces
### BEGIN INIT INFO
# Provides:          check-network-interfaces
# Required-Start:    $local_fs
# X-Start-Before:    ifupdown
# Required-Stop:
# Default-Start:     S
# Default-Stop:
### END INIT INFO
# Reported to BTS as #459872
