# debian/overrides/console-screen.sh
### BEGIN INIT INFO
# Provides:          console-screen
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     S 2 3 4 5
# Default-Stop:      0 1 6
# Should-Start:      $syslog
# Should-Stop:
### END INIT INFO
# Reported to BTS as bug #330229
