# debian/overrides/cpqarrayd
### BEGIN INIT INFO
# Provides:          cpqarrayd
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO
