# debian/overrides/dbus-1
### BEGIN INIT INFO
# Provides:          dbus-1
# Required-Start:    $local_fs $syslog
# Required-Stop:     $local_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO
# Reported to Debian BTS as #337644
