# debian/overrides/discover
### BEGIN INIT INFO
# Provides:          discover
# Required-Start:    mountdevsubfs checkroot $local_fs
# Required-Stop: 
# Default-Start:     S
# Default-Stop:
### END INIT INFO
# Discover v2 reported to BTS as #335322
# Discover v1 reported to BTS as #335323
