# debian/overrides/dmraid
### BEGIN INIT INFO
# Provides:          dmraid
# Required-Start:    mountkernfs udev
# Required-Stop:     umountroot
# X-Start-Before:    checkroot mdadm-raid lvm
# X-Stop-After:      umountfs mdadm-raid lvm
# Default-Start:     S
# Default-Stop:      0 6
### END INIT INFO
# Reported to BTS as #458626
