# debian/overrides/flashybrid
### BEGIN INIT INFO
# Provides:          flashhybrid
# Required-Start:    mountdevsubfs $local_fs
# Required-Stop:     mountdevsubfs $local_fs
# Default-Start:     S 6
# Default-Stop:
### END INIT INFO
