# debian/overrides/gdm
### BEGIN INIT INFO
# Provides:          gdm
# Required-Start:    $syslog $local_fs
# Required-Stop:     $syslog $local_fs
# Should-Start:      xfs
# Should-Stop:       xfs
# Default-Start:     2 3 4 5
# Default-Stop:      S 0 1 6
### END INIT INFO
