# debian/overrides/gpsd
### BEGIN INIT INFO
# Provides:          gpsd
# Required-Start:    $syslog $network
# Required-Stop:     $syslog $network
# Default-Start:     S 1 2 3 4 5
# Default-Stop:      0 6
### END INIT INFO
