# debian/overrides/hdparm
### BEGIN INIT INFO
# Provides:          hdparm
# Required-Start:    mountdevsubfs
# Should-Start:      udev
# Required-Stop:
# Default-Start:     S
# Default-Stop:      0 6
### END INIT INFO
# Reported to BTS as #335309
