# debian/overrides/hibernate
### BEGIN INIT INFO
# Provides:          hibernate
# Required-Start:    checkfs
# Required-Stop: 
# Default-Start:     S
# Default-Stop:
### END INIT INFO
