# debian/overrides/hostname.dhcp
### BEGIN INIT INFO
# Provides:          hostname.dhcp
# Required-Start:    $local_fs hostname
# Required-Stop:
# Should-Start:      $network
# X-Start-Before:    $remote_fs
# Default-Start:     S
# Default-Stop:
### END INIT INFO
# Reported to BTS as #458859
