# debian/overrides/hotplug-net
### BEGIN INIT INFO
# Provides:          hotplug-net
# Required-Start:    ifupdown
# Required-Stop:     ifupdown
# Default-Start:     S
# Default-Stop:      0 6
### END INIT INFO
