# debian/overrides/hplip
### BEGIN INIT INFO
# Provides:          hplip
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Should-Start:      $network
# Should-Stop:       $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO
# Reported to bTS as #427036
