# debian/overrides/hwclock.sh
### BEGIN INIT INFO
# Provides:          hwclock
# Required-Start:    console-screen
# Required-Stop:     console-screen
# Default-Start:     S 1 2 3 4 5
# Default-Stop:      0 6
### END INIT INFO
