# debian/overrides/hwclockfirst.sh
### BEGIN INIT INFO
# Provides:          hwclockfirst
# Required-Start:    
# Required-Stop:     
# Default-Start:     S
# Default-Stop:
### END INIT INFO
# Reported to BTS as bug #330227
