# debian/overrides/ifplugd
### BEGIN INIT INFO
# Provides:          ifplugd
# Required-Start:    $network $syslog
# Required-Stop:     $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO
