# debian/overrides/ifupdown
### BEGIN INIT INFO
# Provides:          ifupdown
# Required-Start:    ifupdown-clean
# Required-Stop:     ifupdown-clean
# Default-Start:     S
# Default-Stop:      0 6
### END INIT INFO
# Reported to BTS as bug #330230
