# debian/overrides/initrd-tools.sh
### BEGIN INIT INFO
# Provides:          initrd-tools
# Required-Start:    bootlogd
# Required-Stop: 
# Default-Start:     S 1 2 3 4 5
# Default-Stop:      0 6
### END INIT INFO
