# debian/overrides/initrd-tools.sh
### BEGIN INIT INFO
# Provides:          initrd-tools
# Required-Start:    mountvirtfs
# Required-Stop: 
# Default-Start:     S
# Default-Stop:
### END INIT INFO
