# debian/overrides/keymap.sh
### BEGIN INIT INFO
# Provides:          keymap
# Required-Start:    mountdevsubfs
# Required-Stop:     
# Default-Start:     S
# Default-Stop:
# Should-Start:
# Should-Stop:
### END INIT INFO
# Reported to BTS as bug #330231
