# debian/overrides/keymap.sh
### BEGIN INIT INFO
# Provides:          keymap
# Required-Start:    mountdevsubfs
# Required-Stop:     
# Default-Start:     S 2 3 4 5
# Default-Stop:      0 1 6
# Should-Start:
# Should-Stop:
### END INIT INFO
# Reported to BTS as bug #330231
