# debian/overrides/libdevmapper1.00
### BEGIN INIT INFO
# Provides:          libdevmapper
# Required-Start:
# Required-Stop:
# Should-Start:      modutils module-init-tools hotplug discover udev devfsd
# Should-Stop:       modutils module-init-tools
# Default-Start:     S
# Default-Stop:
### END INIT INFO
