# debian/overrides/libdevmapper1.02.1
### BEGIN INIT INFO
# Provides:          libdevmapper
# Required-Start:    modutils
# Required-Stop:     modutils
# Should-Start:      hotplug discover udev devfsd
# Default-Start:     S
# Default-Stop:
### END INIT INFO
