# debian/overrides/libnss-ldap
### BEGIN INIT INFO
# Provides:          libnss-ldap
# Required-Start:    mountkernfs
# Required-Stop:
# Default-Start:     S
# Default-Stop:      0 2 3 4 5 6
### END INIT INFO
