# debian/overrides/lpd
### BEGIN INIT INFO
# Provides:          lpd
# Required-Start:    $syslog
# Required-Stop:     $syslog
# Should-Start:      $network
# Should-Stop:       $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO
