# debian/overrides/ltsp-client
### BEGIN INIT INFO
# Provides:          ltsp-client
# Required-Start:    ltsp-client-setup $network $syslog
# Required-Stop:     ltsp-client-setup $network $syslog
# Should-Start:      xdebconfigurator
# Default-Start:     2 3 4 5
# Default-Stop:      S 0 1 6
### END INIT INFO
# Already present in package
