# debian/overrides/ltsp-client-setup
### BEGIN INIT INFO
# Provides:          ltsp-client-setup
# Required-Start:    $local_fs
# Required-Stop:     $local_fs
# Default-Start:     S
# Default-Stop:      0 6
### END INIT INFO
# Already present in package
