# debian/overrides/lvm
### BEGIN INIT INFO
# Provides:          lvm
# Required-Start:    mountdevsubfs
# Required-Stop:
# Should-Start:      udev devfsd raid2 mdadm
# Should-Stop:       reboot halt
# Default-Start:     S
# Default-Stop:      0 6
### END INIT INFO
