# debian/overrides/lvm2
### BEGIN INIT INFO
# Provides:          lvm
# Required-Start:    mountdevsubfs
# Required-Stop:
# Should-Start:      udev devfsd raid2 mdadm cryptdisks-early
# Should-Stop:       reboot halt
# Default-Start:     S
# Default-Stop:      0 6
### END INIT INFO
# Reported to BTS as #426108
