# debian/overrides/makedev
### BEGIN INIT INFO
# Provides:          makedev
# Required-Start:    $local_fs
# Required-Stop:     $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:
### END INIT INFO
