# debian/overrides/mdadm
### BEGIN INIT INFO
# Provides:          mdadm
# Required-Start:    mountvirtfs
# Required-Stop:     mountvirtfs
# Should-Start:      devfsd udev discover hotplug
# Default-Start:     S 1 2 3 4 5
# Default-Stop:      0 6
### END INIT INFO
