# debian/overrides/mdadm-raid
### BEGIN INIT INFO
# Provides:          mdadm-raid
# Required-Start:    modutils
# Required-Stop:     modutils
# Should-Start:      hotplug discover udev devfsd
# Default-Start:     S
# Default-Stop:      0 6
### END INIT INFO
