# debian/overrides/modutils
### BEGIN INIT INFO
# Provides:          modutils
# Required-Start:    hwclockfirst checkroot
# Required-Stop:     hwclockfirst checkroot
# Default-Start:     S 1 2 3 4 5
# Default-Stop:      0 6
### END INIT INFO
