# debian/overrides/mountall.sh
### BEGIN INIT INFO
# Provides:          mountall
# Required-Start:    checkfs
# Required-Stop: 
# Should-Start:      lvm
# Should-Stop:       lvm
# Default-Start:     S
# Default-Stop:
### END INIT INFO
