# debian/overrides/mountnfs.sh
### BEGIN INIT INFO
# Provides:          mountnfs nfs
# Required-Start:    $network $local_fs
# Required-Stop:
# Default-Start:     S
# Default-Stop:
### END INIT INFO
