# debian/overrides/mountnfs.sh
### BEGIN INIT INFO
# Provides:          mountnfs nfs
# Required-Start:    $network $local_fs
# Required-Stop: 
# Default-Start:     S 1 2 3 4 5
# Default-Stop:      0 6
### END INIT INFO
