# debian/overrides/mountall.sh
### BEGIN INIT INFO
# Provides:          mountall
# Required-Start:    checkfs
# Required-Stop: 
# Should-Start:      lvm
# Should-Stop:       lvm
# Default-Start:     S 1 2 3 4 5
# Default-Stop:      0 6
### END INIT INFO
