# debian/overrides/nbd-client
### BEGIN INIT INFO
# Provides:          nbd-client
# Required-Start:    $network
# X-Start-Before:    $remote_fs
# Required-Stop:     $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO
