# debian/overrides/networking
### BEGIN INIT INFO
# Provides:          networking
# Required-Start:    mountvirtfs ifupdown $local_fs
# Required-Stop:     mountvirtfs
# Should-Start:      hotplug discover
# Default-Start:     S 1 2 3 4 5
# Default-Stop:      0 6
### END INIT INFO
