# debian/overrides/networking
### BEGIN INIT INFO
# Provides:          networking
# Required-Start:    mountdevsubfs ifupdown $local_fs
# Default-Start:     S
# Default-Stop:      0 6
### END INIT INFO
