# debian/overrides/nfs-common
### BEGIN INIT INFO
# Provides:          nfs-common
# Required-Start:    $time
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO
