# debian/overrides/nfs-common
### BEGIN INIT INFO
# Provides:          nfs-common
# Required-Start:    $network $portmap
# Required-Stop:     $network $portmap
# Default-Start:     2 3 4 5
# Default-Stop:      S 0 1 6
### END INIT INFO
