# debian/overrides/nscd
### BEGIN INIT INFO
# Provides:          nscd
# Required-Start:    $syslog
# Required-Stop:     $syslog
# Should-Start:      $network slapd $named
# Should-Stop:       $network slapd $named
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO
# Reported to BTS as #335343
