# debian/overrides/nstxcd
### BEGIN INIT INFO
# Provides:          nstxcd
# Required-Start:    $network
# Required-Stop:     $network
# Should-Start:      nstxd
# Should-Stop:       nstxd
# Default-Start:     S 1 2 3 4 5
# Default-Stop:      0 6
### END INIT INFO
