# debian/overrides/nstxd
### BEGIN INIT INFO
# Provides:          nstxd
# Required-Start:    $network
# Required-Stop:     $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO
