# debian/overrides/nstxd
### BEGIN INIT INFO
# Provides:          nstxd
# Required-Start:    $network
# Required-Stop:     $network
# Default-Start:     S 1 2 3 4 5
# Default-Stop:      0 6
### END INIT INFO
