# debian/overrides/openbsd-inetd
### BEGIN INIT INFO
# Provides:          openbsd-inetd
# Required-Start:    $network
# Required-Stop:     $network
# Default-Start:     2 3 4 5
# Default-Stop:      S 0 1 6
### END INIT INFO
