# debian/overrides/pcmcia
### BEGIN INIT INFO
# Provides:          pcmcia
# Default-Start:     1 2 3 4 5
# Default-Stop:      0 6
### END INIT INFO
