# debian/overrides/portmap
### BEGIN INIT INFO
# Provides:          portmap
# Required-Start:    $network
# Required-Stop:     $network
# Should-Start:      $named
# Should-Stop:       $named
# Default-Start:     1 2 3 4 5
# Default-Stop:      S 0 6
### END INIT INFO
