# debian/overrides/postfix
### BEGIN INIT INFO
# Provides:	  postfix
# Required-Start: $network $syslog
# Required-Stop:  $network $syslog
# Should-Start:   spampd amavis postgrey
# Should-Stop:    spampd amavis postgrey
# Default-Start:  2 3 4 5
# Default-Stop:   S 0 1 6
### END INIT INFO
