# debian/overrides/pptpd
### BEGIN INIT INFO
# Provides:          pptpd
# Required-Start:    $syslog
# Required-Stop:     $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      S 0 1 6
### END INIT INFO
