# debian/overrides/puppet
### BEGIN INIT INFO
# Provides:          puppet
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Should-Start:      puppetmaster
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO
# Reported to BTS as #?
