# debian/overrides/racoon
### BEGIN INIT INFO
# Provides:          racoon
# Required-Start:    $remote_fs
# Required-Stop:     $remote_fs
# Default-Start:     S
# Default-Stop:      1
### END INIT INFO
# Reported to BTS as #458488
