# debian/overrides/raid2
### BEGIN INIT INFO
# Provides:          raid2
# Required-Start:    modutils
# Required-Stop:     modutils
# Should-Start:      devfsd udev discover hotplug
# Default-Start:     S 1 2 3 4 5
# Default-Stop:      0 6
### END INIT INFO
