# debian/overrides/readahead
### BEGIN INIT INFO
# Provides:          readahead
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     S
# Default-Stop:      0 6
### END INIT INFO
