# debian/overrides/reboot
### BEGIN INIT INFO
# Provides:       reboot
# Required-Start: umountfs $network
# Required-Stop:
# Should-Start:   lvm raid2
# Default-Start:  6
# Default-Stop:
### END INIT INFO
# Already present in package
