# debian/overrides/rsync
### BEGIN INIT INFO
# Provides:          rsyncd
# Provides:          rsyncd
# Required-Start:    $network
# Required-Stop:     $syslog
# Should-Start:      $named $syslog
# Should-Stop:       $syslog network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO
