# debian/overrides/samba
### BEGIN INIT INFO
# Provides:          samba
# Required-Start:    $network
# Required-Stop:     $network
# Should-Start:      $named
# Should-Stop:       $named
# Default-Start:     2 3 4 5
# Default-Stop:      S 0 1 6
### END INIT INFO
