# debian/overrides/screen-cleanup
### BEGIN INIT INFO
# Provides:          screen-cleanup
# Required-Start:    $remote_fs
# Required-Stop:
# Default-Start:     S
# Default-Stop:
### END INIT INFO
