# debian/overrides/scsitools-pre.sh
### BEGIN INIT INFO
# Provides:          scsitool-pre
# Required-Start:    mountdevsubfs
# Required-Stop:
# X-Start-Before:    checkroot
# Default-Start:     S
# Default-Stop:
### END INIT INFO
# Reported to BTS as #458649
