# debian/overrides/sendsigs
### BEGIN INIT INFO
# Provides:          sendsigs
# Required-Start:    
# Required-Stop: 
# Default-Start:      6
# Default-Stop:       
# Short-Description: 
# Description: 
### END INIT INFO
# Already present in package
