# debian/overrides/slashem-common
### BEGIN INIT INFO
# Provides:          slashem-common
# Required-Start:    $remote_fs
# Required-Stop:     $remote_fs
# Default-Start:     S
# Default-Stop:
### END INIT INFO
# Reported to BTS as #465433
