# debian/overrides/start-wlan
### BEGIN INIT INFO
# Provides:          start-wlan
# Required-Start:    $network
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      S 1 0 6
### END INIT INFO
