# debian/overrides/sysfsutils
### BEGIN INIT INFO
# Provides:          sysfsutils
# Required-Start:    $local_fs $remote_fs
# Required-Stop: 
# Default-Start:     2 3 4 5
# Default-Stop:
### END INIT INFO
# Reported to BTS as #427033
