# debian/overrides/sysklogd
### BEGIN INIT INFO
# Provides:          sysklogd
# Required-Start:    $local_fs $time
# Required-Stop:     $local_fs $time
# Should-Start:      $network $named
# Should-Stop:       $network $named
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO
# Reported to BTS as bug #324679
