# debian/overrides/timidity
### BEGIN INIT INFO
# Provides:          timidity
# Required-Start:    $remote_fs
# Required-Stop:     $remote_fs
# Should-Start:      alsa-utils
# Should-Stop:       alsa-utils
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO
