# debian/overrides/udev
### BEGIN INIT INFO
# Provides:          udev
# Required-Start:    mountkernfs
# Required-Stop:     
# Default-Start:     S
# Default-Stop:
### END INIT INFO
# Reported to BTS as #335310
