# debian/overrides/udev-mtab
### BEGIN INIT INFO
# Provides:          udev-mtab
# Required-Start:    $local_fs
# Required-Stop:     $local_fs
# Default-Start:     S
# Default-Stop:      0 1 6
### END INIT INFO
# Reported to BTS as #335310
