# debian/overrides/udev-mtab
### BEGIN INIT INFO
# Provides:          udev-mtab
# Required-Start:    $local_fs
# Required-Stop:     $local_fs
# Default-Start:     S 1 2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO
