# debian/overrides/umountnfs.sh
### BEGIN INIT INFO
# Provides:          umountnfs
# Required-Start:    
# Required-Stop: 
# Default-Start:      6
# Default-Stop:       
### END INIT INFO
