# debian/overrides/umountnfs.sh
### BEGIN INIT INFO
# Provides:       umountnfs
# Required-Start: sendsigs
# Required-Stop:
# Default-Start:  6
# Default-Stop:
### END INIT INFO
# Already present in package
