# debian/overrides/umountfs
### BEGIN INIT INFO
# Provides:       umountfs
# Required-Start: umountnfs
# Required-Stop:
# Default-Start:  6
# Default-Stop:
### END INIT INFO
# Already present in package
