# debian/overrides/urandom
### BEGIN INIT INFO
# Provides:          urandom
# Required-Start:    $local_fs
# Required-Stop:     $local_fs
# Default-Start:     S 1 2 3 4 5
# Default-Stop:      0 6
### END INIT INFO
