# debian/overrides/whereami
### BEGIN INIT INFO
# Provides:          whereami
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      S 0 1 6
### END INIT INFO
