# debian/overrides/whitelister
### BEGIN INIT INFO
# Provides:          whitelister
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# X-Start-Before:    postfix
# X-Stop-After:      postfix
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO
# Reported to the BTS as #466659
