# debian/overrides/x11-common
### BEGIN INIT INFO
# Provides:          x11-common
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Should-Start:      $named
# Should-Stop:       $named
# Default-Start:     S
# Default-Stop:
### END INIT INFO
