# debian/overrides/xserver-xorg
### BEGIN INIT INFO
# Provides:          xserver-xorg
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     S 1 2 3 4 5
# Default-Stop:      0 6
### END INIT INFO
